# var/dpkg/info/sasl2-bin.list
/.
/etc
/etc/default
/etc/default/saslauthd
/etc/init.d
/etc/init.d/saslauthd
/usr
/usr/bin
/usr/bin/gen-auth
/usr/bin/sasl-sample-client
/usr/bin/saslfinger
/usr/lib
/usr/lib/sasl2
/usr/lib/sasl2/berkeley_db.txt
/usr/sbin
/usr/sbin/sasl-sample-server
/usr/sbin/saslauthd
/usr/sbin/sasldbconverter2
/usr/sbin/sasldblistusers2
/usr/sbin/saslpasswd2
/usr/sbin/saslpluginviewer
/usr/sbin/testsaslauthd
/usr/share
/usr/share/doc
/usr/share/doc/sasl2-bin
/usr/share/doc/sasl2-bin/README.Debian.gz
/usr/share/doc/sasl2-bin/TODO.Debian
/usr/share/doc/sasl2-bin/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/sasl2-bin
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/gen-auth.1.gz
/usr/share/man/man1/saslfinger.1.gz
/usr/share/man/man8
/usr/share/man/man8/sasl-sample-client.8.gz
/usr/share/man/man8/sasl-sample-server.8.gz
/usr/share/man/man8/saslauthd.8.gz
/usr/share/man/man8/sasldblistusers2.8.gz
/usr/share/man/man8/saslpasswd2.8.gz
/usr/share/man/man8/saslpluginviewer.8.gz
/usr/share/man/man8/testsaslauthd.8.gz
/usr/share/doc/sasl2-bin/NEWS.Debian.gz
/usr/share/doc/sasl2-bin/changelog.Debian.gz
