# var/dpkg/info/tmpreaper.list
/.
/etc
/etc/cron.daily
/etc/cron.daily/tmpreaper
/etc/tmpreaper.conf
/usr
/usr/sbin
/usr/sbin/tmpreaper
/usr/share
/usr/share/doc
/usr/share/doc/tmpreaper
/usr/share/doc/tmpreaper/README.security.gz
/usr/share/doc/tmpreaper/changelog.gz
/usr/share/doc/tmpreaper/copyright
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/tmpreaper.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/tmpreaper.8.gz
