# var/dpkg/info/wide-dhcpv6-client.conffiles
/etc/init.d/wide-dhcpv6-client
/etc/wide-dhcpv6/dhcp6c-ifupdown
/etc/wide-dhcpv6/dhcp6c-script
/etc/wide-dhcpv6/dhcp6c.conf
